# root/etc/s6-overlay/s6-rc.d/tinyproxy/run
#!/usr/bin/with-contenv bash
if [[ $HTTP -eq 0 ]]; then
    tail -f /dev/null
    exit
fi
tinyproxy -d -c /etc/tinyproxy/tinyproxy.conf
