# root/etc/s6-overlay/s6-rc.d/wg/run
#!/usr/bin/with-contenv bash
resolvconf -u
wg-quick up wg0
ip rule add from $(hostname -i) lookup main
wg
